# selenium/flask.conf
[program:flaskapp]
priority=30
user=seluser
command=/opt/bin/start-flask-server.sh
autostart=true
autorestart=false
startsecs=0
startretries=0
